# systemd/dhclient-wlan0.service
[Unit]
Description=Run dhclient on wlan0 after Wi-Fi connects
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/dhclient wlan0
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
